# otaserver.conf
start on filesystem and started networking
respawn
chdir /node-esp-OTAserver  #deployment directory
env NODE_ENV=production
env PORT=80
exec /usr/local/bin/node ./bin/www